# systemd/denmarkapi-translate.service
[Unit]
Description=DenmarkAPI smiley report English translation (gpt-oss-20b via local vLLM)
After=network-online.target docker.service

[Service]
Type=simple
User=arian
WorkingDirectory=/home/arian/Programming/DenmarkAPI
# -u so the journal shows progress live. No --concurrency: it follows the OVERLAY slider,
# because it is the same workload against the same vLLM and shares the same paragraph cache.
# Best run AFTER/BEHIND the overlay, where it is nearly free (~250 reports/s, 0 LLM calls).
ExecStart=/home/arian/Programming/DenmarkAPI/.venv/bin/python -u -m denmarkapi.smiley.translate --watch
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
